# deployment/supervisor/vendor.conf
[program:vendor]
command=/nfs/soft/www/apps/vendor_uploads/envs/development/bin/gunicorn -b 0.0.0.0:5001 -w 4 vendor:app --timeout 300
directory=/nfs/soft/www/apps/vendor_uploads/dist/vendor-upload
user=www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true